# hostsupdate/hostsupdate.service
[Unit]
Description='hosts.d update'

[Service]
Type=oneshot
User=dnsmasq
Group=dnsmasq
ExecStart=/usr/bin/hostsupdate.py
