# hostsupdate/hostsupdate.timer
[Unit]
Description='Weekly hosts.d update'

[Timer]
OnCalendar=weekly
Persistent=true
AccuracySec=1h

[Install]
WantedBy=timers.target
