# domainblacklist-git/domainblacklist.timer
[Unit]
Description='Weekly domain blacklists update'

[Timer]
OnCalendar=weekly
Persistent=true
AccuracySec=1h

[Install]
WantedBy=timers.target
