# domainblacklist/domainblacklist.service
[Unit]
Description='Domain blacklists update'

[Service]
Type=oneshot
User=nobody
Group=nobody
ExecStart=/usr/bin/domainblacklist
